# provisioning/remote_files/erlangelist-site.service
[Unit]
Description=The Erlangelist site

[Service]
ExecStart=/opt/erlangelist/erlangelist-site.sh startf
ExecStop=/opt/erlangelist/erlangelist-site.sh stop
Restart=always
RestartSec=1s

[Install]
WantedBy=multi-user.target
